# systemd/conky.timer
[Unit]
Description=Delayed start for conky

[Timer]
OnActiveSec=5
Unit=conky.service

[Install]
WantedBy=default.target
